# systemd/user/disable-suspend-on-lid-close.service
[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/libexec/gnome-tweak-tool-lid-inhibitor
[Install]
WantedBy=default.target

